# deploy/autoupdate.service

[Unit]
Description=NiketBot auto-update from git
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/niket_ai_trader
ExecStart=/usr/bin/git pull --rebase --autostash
ExecStartPost=%h/niket_ai_trader/deploy/install_or_update.sh
